# .config/systemd/user/grive.service
[Unit]
Description=Grive service
[Service]
ExecStart=/home/deleuze/.bin/grive-sync
[Install]
WantedBy=default.target
