# modules/rclocal_config/files/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
rfkill block bluetooth
exit 0
